# supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2 -D FOREGROUND
redirect_stderr=true
killasgroup=true
stopasgroup=true

[program:cron]
command=cron -f
